# launchd/com.nanoclaw.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.web</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{NODE_PATH}}</string>
        <string>{{PROJECT_ROOT}}/web/node_modules/next/dist/bin/next</string>
        <string>start</string>
        <string>-p</string>
        <string>3030</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}/web</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{HOME}}/.bun/bin:{{HOME}}/.orbstack/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:{{HOME}}/.local/bin</string>
        <key>HOME</key>
        <string>{{HOME}}</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>NANOCLAW_HTTP_HOST</key>
        <string>127.0.0.1</string>
        <key>NANOCLAW_HTTP_PORT</key>
        <string>7878</string>
    </dict>
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/nanoclaw.web.log</string>
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/nanoclaw.web.error.log</string>
</dict>
</plist>
